# Server/rot13.socket
[Unit]
Description=Socket
PartOf=rot13.service

[Socket]
ListenStream=0.0.0.0:8080
Accept=yes
FreeBind=yes

[Install]
WantedBy=sockets.target